# src/config/supervisor/scheduler.conf
[program:scheduler]
process_name=%(program_name)s_%(process_num)02d
# Running artisan command in the background every one minute (60 seconds)
command=/bin/sh -c 'while true; do date && /usr/local/bin/php /app/artisan schedule:run; sleep 60; done'
autostart=true
autorestart=true
user=root
numprocs=1
# See http://supervisord.org/configuration.html#program-x-section-settings
stopwaitsecs=120
# INT => 2
# TERM => 15
# QUIT => 3
# See https://de.wikipedia.org/wiki/Signal_(Unix)
# See http://supervisord.org/configuration.html?highlight=stopsignal
stopsignal=QUIT
# redirect logs to stdout
# see http://veithen.io/2015/01/08/supervisord-redirecting-stdout.html
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
